# battle-buddy-law.service
[Unit]
Description=Battle Buddy Listener — Travis County Law
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/battle_buddy
EnvironmentFile=/home/pi/battle_buddy/config.env
ExecStart=/usr/bin/python3 /home/pi/battle_buddy/battle_buddy_listener_v1.2.py --stream law
Restart=always
RestartSec=15
StandardOutput=append:/home/pi/battle_buddy/logs/listener_law.log
StandardError=append:/home/pi/battle_buddy/logs/listener_law.log

[Install]
WantedBy=multi-user.target
